# root/etc/services.d/rtorrent/run
#!/usr/bin/execlineb -P

with-contenv
s6-setuidgid nginx
rtorrent -n -o import=/etc/rtorrent/rtorrent.rc
