# root/etc/services.d/nginx/run
#!/usr/bin/execlineb -P

with-contenv
s6-setuidgid nginx
nginx -g "daemon off;"
